# recipes-phosphor/mctp/files/mctp-setup-i2c.service
[Unit]
Description=MCTP I2C transport setup for QEMU AST2600 (i2c-1)
# Start after mctpd is active; stop if mctpd stops.
After=mctpd.service
Before=pldmd.service
BindTo=mctpd.service

[Service]
Type=oneshot
RemainAfterExit=yes

# The mctp-i2c driver (CONFIG_MCTP_TRANSPORT_I2C=y, built-in) checks for the
# "mctp-controller" OF property on the adapter to create the mctpi2cN netdev.
# The DTS patch adds that property to i2c1, but the kernel notifier only fires
# for adapters registered AFTER the mctp-i2c init runs (init order is i2c before
# net, so aspeed adapters are already registered by then). We must trigger probe
# manually via new_device so its bus_for_each_dev scan finds i2c-1.
# NOTE: the i2c_device_id name is "mctp-i2c-interface", NOT "mctp-i2c".
ExecStart=-/bin/sh -c 'echo "mctp-i2c-interface 0x0f" > /sys/bus/i2c/devices/i2c-1/new_device'

# Wait for the probe to complete and mctpi2c1 netdev to appear (asynchronous).
ExecStart=/bin/sh -c 'for i in $(seq 1 10); do [ -d /sys/class/net/mctpi2c1 ] && exit 0; sleep 1; done; echo "mctpi2c1 did not appear" >&2; exit 1'

# Bring the MCTP I2C netdev up.
# NOTE: BMC EID 8 is intentionally NOT added here.
# mctp-lo-setup.service adds EID 10 first (making it addrs[0]) and then EID 8
# (addrs[1]). The kernel always uses addrs[0] as the source EID for outgoing
# packets (mctp_local_output). EID 10 must be addrs[0] so that loopback
# responses from pldm-terminus carry src=10, matching pldmd's waiting key
# {local=10, peer=10}. If EID 8 were addrs[0], responses would carry src=8
# and mctp_key_match would fail.
ExecStart=/usr/bin/mctp link set mctpi2c1 up

[Install]
WantedBy=mctpd.service
