# recipes-phosphor/mctp/files/mctp-setup-i2c.service
[Unit]
Description=MCTP I2C transport setup for QEMU AST2600 (i2c-1)
# Start after mctpd is active; stop if mctpd stops.
After=mctpd.service
Before=pldmd.service
BindTo=mctpd.service

[Service]
Type=oneshot
RemainAfterExit=yes

# The mctp-i2c driver (CONFIG_MCTP_TRANSPORT_I2C=y, built-in) checks for the
# "mctp-controller" OF property on the adapter to create the mctpi2cN netdev.
# The DTS patch adds that property to i2c1, but the kernel notifier only fires
# for adapters registered AFTER the mctp-i2c init runs (init order is i2c before
# net, so aspeed adapters are already registered by then). We must trigger probe
# manually via new_device so its bus_for_each_dev scan finds i2c-1.
# NOTE: the i2c_device_id name is "mctp-i2c-interface", NOT "mctp-i2c".
#
# NOTE: this instantiates the BMC's OWN local port on the bus at
# BMC_I2C_SLAVE_ADDR (0x10), NOT the terminus's address (0x0f). The QEMU
# mctp-i2c-endpoint device (the remote terminus) sits at 0x0f and sends its
# responses via a master-write targeting 0x10; the local kernel client must
# listen there. Using 0x0f here collides with the terminus's own address
# on the emulated bus and no response can ever be received.
ExecStart=-/bin/sh -c 'echo "mctp-i2c-interface 0x10" > /sys/bus/i2c/devices/i2c-1/new_device'

# Wait for the probe to complete and mctpi2c1 netdev to appear (asynchronous).
ExecStart=/bin/sh -c 'for i in $(seq 1 10); do [ -d /sys/class/net/mctpi2c1 ] && exit 0; sleep 1; done; echo "mctpi2c1 did not appear" >&2; exit 1'

# Bring the MCTP I2C netdev up and add the BMC's own EID.
ExecStart=/usr/bin/mctp link set mctpi2c1 up
ExecStart=/usr/bin/mctp addr add 8 dev mctpi2c1

[Install]
WantedBy=mctpd.service
