# recipes-phosphor/pldm-terminus/files/mctp-lo-setup.service
[Unit]
Description=MCTP local EID setup for PLDM terminus testing
# Run after mctpi2c1 is up. EID 8 is NOT added by mctp-setup-i2c anymore;
# we add EID 10 first here so it becomes addrs[0] on mctpi2c1.
After=mctp-setup-i2c.service
Wants=mctp-setup-i2c.service
# Must complete before pldmd starts so EID 10 LOCAL route exists when pldmd probes.
Before=pldmd.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Add EID 10 FIRST so it becomes addrs[0] on mctpi2c1.
# The kernel (mctp_local_output) always uses addrs[0] as the source EID for
# outgoing packets. With EID 10 as addrs[0], pldmd's requests carry src=10
# and its waiting key becomes {local=10, peer=10}. pldm-terminus responses
# also carry src=10, so mctp_key_match succeeds and pldmd receives the reply.
ExecStart=/usr/bin/mctp addr add 10 dev mctpi2c1

# Add EID 8 second (addrs[1]) so the RTN_LOCAL route for EID 8 exists,
# allowing pldm-terminus to route its response to pldmd via mctp_dst_input.
ExecStart=/usr/bin/mctp addr add 8 dev mctpi2c1

[Install]
WantedBy=multi-user.target
