# recipes-phosphor/pldm-terminus/files/pldm-terminus.service
[Unit]
Description=PLDM base terminus responder (EID 10, MCTP loopback)
After=mctp-lo-setup.service
BindTo=mctp-lo-setup.service
# Must be running before pldmd's initial GetPLDMVersion probe.
Before=pldmd.service

[Service]
Type=simple
ExecStart=/usr/bin/pldm-terminus
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
